# systemd/hs-data-api-docker-firecrawl-streamer.service
[Unit]
Description=Refresh HSGuru streamer decks via Firecrawl in Docker
Requires=hs-data-api-docker.service
After=hs-data-api-docker.service network-online.target

[Service]
Type=oneshot
WorkingDirectory=/srv/hs-data-api
ExecStart=/usr/bin/docker compose -f /srv/hs-data-api/docker-compose.yml run --rm api python /app/scripts/firecrawl-streamer-decks.py
