# systemd/hs-data-api-docker-refresh-fun-decks-standard.service
[Unit]
Description=Refresh HSGuru fun/off-meta decks with Standard focus in Docker
Requires=hs-data-api-docker.service
After=hs-data-api-docker.service network-online.target

[Service]
Type=oneshot
WorkingDirectory=/srv/hs-data-api
ExecCondition=/usr/bin/docker compose -f /srv/hs-data-api/docker-compose.yml run --rm api python -m app.cli scheduled-check --source hsguru_fun_decks
ExecStart=/usr/bin/docker compose -f /srv/hs-data-api/docker-compose.yml run --rm api python -m app.cli refresh-fun-decks --scheduled --format standard
