# systemd/hs-data-api-docker-refresh-hsreplay-archetypes.timer
[Unit]
Description=Refresh HSReplay Standard archetype SQLite database twice per week

[Timer]
OnCalendar=Mon,Thu *-*-* 03:20:00 Europe/Warsaw
AccuracySec=10min
Persistent=true

[Install]
WantedBy=timers.target
