# systemd/hs-data-api-docker-refresh-vicious-syndicate.service
[Unit]
Description=Refresh Vicious Syndicate Live and latest available radars in Docker
Requires=hs-data-api-docker.service
After=hs-data-api-docker.service network-online.target

[Service]
Type=oneshot
WorkingDirectory=/srv/hs-data-api
ExecStart=/usr/bin/docker compose -f /srv/hs-data-api/docker-compose.yml run --rm api python -m app.cli refresh --source vicious_syndicate_live_beta --source vicious_syndicate_radars
