# systemd/hs-data-api-refresh-protected.service
[Unit]
Description=Refresh HSGuru / protected browser sources (night tier)
After=network-online.target hs-scrape-proxy.service hs-flaresolverr.service
Wants=network-online.target hs-scrape-proxy.service hs-flaresolverr.service

[Service]
Type=oneshot
WorkingDirectory=/opt/hs-data-api
EnvironmentFile=/etc/hs-data-api.env
Environment=DISPLAY=:99
ExecStartPre=/opt/hs-data-api/scripts/ensure-flaresolverr.sh
ExecStartPre=/opt/hs-data-api/scripts/ensure-xvfb.sh
ExecStartPre=/opt/hs-data-api/venv/bin/python -m app.cli preflight --strict
ExecStart=/opt/hs-data-api/venv/bin/python -m app.cli refresh --tier browser_protected
