# systemd/hs-data-api-docker-rebuild-hsreplay-index.timer
[Unit]
Description=Rebuild the derived HSReplay index daily after API refresh

[Timer]
OnCalendar=*-*-* 19:15:00 Europe/Warsaw
AccuracySec=10min
Persistent=true

[Install]
WantedBy=timers.target
