# systemd/hs-data-api-docker-refresh-hsguru-meta-matrix.timer
[Unit]
Description=Refresh the unified HSGuru meta matrix every six hours

[Timer]
OnCalendar=*-*-* 00:00:00 Europe/Warsaw
OnCalendar=*-*-* 06:00:00 Europe/Warsaw
OnCalendar=*-*-* 12:00:00 Europe/Warsaw
OnCalendar=*-*-* 18:00:00 Europe/Warsaw
AccuracySec=5min
Persistent=true

[Install]
WantedBy=timers.target
