# systemd/hs-data-api-docker-refresh-bg-hero-details.timer
[Unit]
Description=Refresh HSReplay Battlegrounds hero detail cache twice weekly

[Timer]
OnCalendar=Mon,Thu *-*-* 04:35:00 Europe/Warsaw
AccuracySec=10min
Persistent=true

[Install]
WantedBy=timers.target
