# systemd/hs-data-api-firecrawl-streamer.service
[Unit]
Description=Refresh HSGuru streamer decks via Firecrawl
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/hs-data-api
EnvironmentFile=/etc/hs-data-api.env
ExecCondition=/opt/hs-data-api/venv/bin/python -m app.cli scheduled-check --source hsguru_streamer_decks_legend_1000
ExecStart=/opt/hs-data-api/venv/bin/python /opt/hs-data-api/scripts/firecrawl-streamer-decks.py
ExecStartPost=-/usr/bin/chown -R hs-data-api:hs-data-api /var/lib/hs-data-api

UMask=0022
NoNewPrivileges=true
